# keep-active.timer
[Unit]
Description=Keep active timer

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target
