# keep-active.timer
[Unit]
Description=Keep active timer

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri *-*-* 9..17:*:00

[Install]
WantedBy=timers.target
