# keep-active.service
[Unit]
Description=Keep active job

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/bin/keep-active.sh
